# data/base/units/etcd3-defragmentation.timer
[Unit]
Description=Execute etcd3-defragmentation every day at 3.30AM UTC
[Timer]
OnCalendar=*-*-* 03:30:00 UTC
[Install]
WantedBy=multi-user.target
